# core/pgbackrest/pgbackrest-incr.timer
[Unit]
Description=Run pgBackRest incremental backup hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
